# conf/etc/sv/test/run
#!/bin/sh -e
exec 2>&1
exec chpst -u foo /opt/test/test.sh
